# Docker/supervisor.conf
[supervisord]
nodaemon=true

[program:jackwui]
command=node /opt/jackwui/bin/jackwui
stdout_logfile_maxbytes = 0
redirect_stderr=true
stdout_logfile=/dev/sdtou
